# ev3dev/robot.service
[Unit]
Description=gRPC server for the robot 

[Service]
User=root
ExecStart=/home/robot/main
Restart=always

[Install]
WantedBy=multi-user.target
